# configs/supervisord.d/dci-lab.ini
[program:dci-lab]
command=/usr/bin/python3 /opt/dci-provisioner/api/lab.py
