# configs/supervisord.d/dci-provisioner.ini
[program:dci-provisioner]
command=/usr/bin/python3 /opt/dci-provisioner/api/provisioner.py
