# configs/supervisord.d/hooktftp.ini
[program:hooktftp]
command=/opt/dci-provisioner/hooktftp/hooktftp -v
user=root
