# configs/supervisord.d/hooktftp.ini
[program:hooktftp]
command=/opt/dci-provisioner/hooktftp/hooktftp -v
user=root
stdout_logfile=/var/log/supervisord/hooktftp.log
redirect_stderr=true
