# configs/supervisord.d/rq_worker.ini
[program:redis]
command=/opt/dci-provisioner/run_worker.sh
user=root
stdout_logfile=/var/log/supervisord/rq_worker.log
redirect_stderr=true
