# configs/supervisord.d/rq_worker.ini
[program:redis]
command=/opt/dci-provisioner/run_worker.sh
user=root
